# config/supervisord/supervisord.conf
[program:rsyslogd]
command = /usr/sbin/rsyslogd -n
autostart=true
autorestart=true

[program:cron]
command = /usr/sbin/cron -f -L 15
autostart=true
autorestart=true

[program:sensu-server]
priority=20
directory=/tmp
command=/opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -v -l /var/log/sensu/sensu-server.log
user=sensu
startsecs=5
autostart=true
autorestart=true

[program:sensu-api]
priority=30
directory=/tmp
command=/opt/sensu/bin/sensu-api -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -v -l /var/log/sensu/sensu-api.log
user=sensu
startsecs=5
autostart=true
autorestart=true

[program:sensu-client]
priority=40
directory=/tmp
command=/opt/sensu/bin/sensu-client -c /etc/sensu/config.json -d /etc/sensu/conf.d -e /etc/sensu/extensions -v -l /var/log/sensu/sensu-client.log
user=sensu
autostart=true
autorestart=true

[program:uchiwa]
priority=50
directory=/tmp
command=/opt/uchiwa/bin/uchiwa -c /etc/sensu/uchiwa.json -d /etc/sensu/dashboard.d -p /opt/uchiwa/src/public
user=uchiwa
autostart=true
autorestart=true
